# scripts/systemd/sentinel-scheduler.service
[Unit]
Description=Sentinel Macro Scheduler Daemon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sentinel/sentiment
ExecStart=/home/ubuntu/sentinel/sentiment/venv/bin/python scripts/macro_scheduler_cli.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/sentinel/sentiment/logs/macro_pipeline.log
StandardError=append:/home/ubuntu/sentinel/sentiment/logs/macro_pipeline.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
