# scripts/systemd/sentinel-api.service
[Unit]
Description=Sentinel IntentCore FastAPI Webhook Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sentinel/sentiment
ExecStart=/home/ubuntu/sentinel/sentiment/venv/bin/uvicorn functions.api.server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
