# scripts/systemd/sentinel-orchestrator.service
[Unit]
Description=Sentinel Portfolio Orchestrator Daemon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sentinel/sentiment
ExecStart=/home/ubuntu/sentinel/sentiment/venv/bin/python scripts/sentinel_orchestrator.py --background
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
